# ansible_hu/additions/gmond_conf/gmond.service.d/override.conf
[Service]
Type=exec
ExecStart=
ExecStart=/usr/sbin/gmond -f
TimeoutStartSec=20
Restart=on-failure
RestartSec=10
